# files/etc/crontabs/root
30 23 * * 5 sleep 5 && touch /etc/banner && reboot
1 1 * * * /usr/share/jd-dailybonus/newapp.sh -w
1 1 * * * /usr/share/jd-dailybonus/newapp.sh -u
