# files/etc/crontabs/root
30 23 * * 5 sleep 5 && touch /etc/banner && reboot
5 1 * * * sleep 18s; /usr/share/jd-dailybonus/newapp.sh -w
1 1 * * * /usr/share/jd-dailybonus/newapp.sh -u
0 2 * * * /usr/share/shadowsocksr/ssrplusupdate.sh >>/tmp/ssrplus.log 2>&1
